# deploy/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

# ============================================
# FastAPI 后端服务
# ============================================
[program:backend]
command=bash -c "source venv/bin/activate && python run.py --host 0.0.0.0 --port 8000"
directory=/root/autodl-tmp/qqchat-enhanced/backend
autostart=true
autorestart=true
startsecs=5
startretries=5
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/tmp/backend-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/tmp/backend-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=USE_POSTGRESQL="true",PG_HOST="localhost",PG_PORT="5432",PG_USER="qqassistant",PG_PASSWORD="qqassistant2026",PG_DATABASE="qqassistant",REDIS_URL="redis://localhost:6379/0",VLLM_ENABLED="true"

# ============================================
# Next.js 前端服务
# ============================================
[program:frontend]
command=pnpm dev
directory=/root/autodl-tmp/qqchat-enhanced
autostart=true
autorestart=true
startsecs=10
startretries=3
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/tmp/frontend-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/tmp/frontend-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=PORT="5000"

# ============================================
# NoneBot2 QQ Bot
# ============================================
[program:qqbot]
command=bash -c "source venv/bin/activate && python -m nb run"
directory=/root/autodl-tmp/qqchat-enhanced/backend
autostart=false
autorestart=true
startsecs=5
startretries=5
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/tmp/qqbot-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/tmp/qqbot-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
